# supervisord.conf
[supervisord]
nodaemon=true

[program:jira]
command=/entrypoint
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
